# deploy/systemd/tothemoon-api.service
[Unit]
Description=ToTheMoon3 API (FastAPI/Uvicorn)
After=network.target

[Service]
User=%i
Group=%i
WorkingDirectory=%h/tothemoon3
Environment=DB_DSN=postgresql+asyncpg://postgres:postgres@localhost:5432/tothemoon2
Environment=LOG_LEVEL=INFO
Environment=PUMPPORTAL_ENABLED=false
ExecStart=%h/tothemoon3/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

