# deploy/systemd/tothemoon2-backend.service
[Unit]
Description=ToTheMoon2 Backend API
After=network.target postgresql.service redis-server.service

[Service]
Type=simple
WorkingDirectory=/opt/tothemoon2/backend
EnvironmentFile=/opt/tothemoon2/.env
ExecStart=/opt/tothemoon2/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

