# deploy/systemd/tothemoon2-celery-worker.service
[Unit]
Description=ToTheMoon2 Celery Worker
After=network.target redis-server.service postgresql.service

[Service]
Type=simple
WorkingDirectory=/opt/tothemoon2
EnvironmentFile=/opt/tothemoon2/.env
ExecStart=/opt/tothemoon2/.venv/bin/celery -A app.core.celery_app worker --loglevel=info --concurrency=1 --prefetch-multiplier=1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

